# playbooks/templates/victoriametrics.service.j2
[Unit]
Description=VictoriaMetrics
After=network.target

[Service]
User=root
Group=root
ExecStart=/opt/victoriametrics/bin/victoria-metrics-prod \
  -storageDataPath=/var/lib/victoria-metrics-data \
  -retentionPeriod=1 \
  -httpListenAddr=:8428
Restart=always

[Install]
WantedBy=multi-user.target
